# net-analyzer/nagios-nrpe/files/nrpe
#!/sbin/runscript

opts="${opts} reload"

depend() {
	need net
}

start() {
	ebegin "Starting nrpe"
	start-stop-daemon --start --quiet --name nrpe \
	--startas /usr/nagios/bin/nrpe \
		-c nagios:nagios \
		-- -c /etc/nagios/nrpe.cfg \
		--daemon
	eend $? "Failed to Start nrpe"
}

stop() {
	ebegin "Stopping nrpe"
	start-stop-daemon --stop --quiet -n nrpe
	eend $? "Failed to Stop nrpe"
}

reload() {
	ebegin "Reloading nrpe"
	kill -HUP `pgrep nrpe | head -1`
	eend $? "Failed to reload nrpe"
}

restart() {
	ebegin "Restarting nrpe"
	svc_stop
	svc_start
	eend $? "Failed to Restart nrpe"
}
