# app-emulation/xen-tools/files/xen-watchdog.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/sbin/xenwatchdogd"
command_args="30 15"
