# net-misc/rtsp-simple-server/files/init.d/rtsp-simple-server
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="rtsp-simple-server daemon"
description="ready-to-use RTSP server and proxy"
command=/usr/bin/rtsp-simple-server
command_args="/etc/rtsp-simple-server.yml"
command_background="true"
command_user="nobody:nobody"
pidfile="/run/rtsp-simple-server.pid"

depend() {
    need net
}
